# monitoring-test.service
[Service]
Type=oneshot
ExecStart=/usr/local/bin/monitor-test.sh
StandardOutput=journal
StandardError=journal
